# krator-os/config/systemd/krator.service
[Unit]
Description=Krator AI assistant daemon
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m krator.agents.assistant
Restart=on-failure

[Install]
WantedBy=multi-user.target
