# krator-os/config/systemd/krator.service
[Unit]
Description=Krator AI assistant daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/krator_daemon.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
